# roles/signalk-server/templates/etc/systemd/system/signalk.socket.j2
[Socket]
ListenStream={{ externalport }}


[Install]
WantedBy=sockets.target

